# l2f-daily.service
[Unit]
Description=Literature to facts daily

[Service]
Type=simple 
ExecStart=/home/savchenk/.pyenv/shims/l2f-tools daily

[Install]
WantedBy=default.target
